# install.ausias-default-settings/usr/share/ausias-cc-default-settings/systemd/datos-.direccion.mount
Description=Direccion
After=network.target

[Mount]
What=lildevil.ausias.lan:/srv/userdata/direccion
Where=/datos/.direccion
Type=nfs
Options=_netdev,auto

[Install]
WantedBy=multi-user.target
