# install.ausias-default-settings/usr/share/ausias-cc-default-settings/systemd/datos-.jefatura.mount
Description=Jefatura
After=network.target

[Mount]
What=lildevil.ausias.lan:/srv/userdata/jefatura
Where=/datos/.jefatura
Type=nfs
Options=_netdev,auto

[Install]
WantedBy=multi-user.target
