# install.ausias-mirror/usr/share/soulmate/profile.d/soulmate-extra-vars.sh
DEBIAN_FORCEMIRROR="http://xwarrior/mirror/debian/"
SM_FORCEMIRROR="http://xwarrior/mirror/soulmate/"
DEBIAN_FORCESECURITY="http://xwarrior/mirror/debian-security/"

