# debian/aurex-lpadmin.prerm
#!/bin/sh
set -e
invoke-rc.d aurex-lpadmin stop || true

