# debian/aurex-lpadmin.upstart
# update aurex-lpadmin printers
#

description     "update aurex-lpadmin printers"

start on started cups

task
exec aurex-lpadmin update

